# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
user=runner

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:worker_submit]
command=python3 -u worker_submit.py
directory=/home/runner/workspace
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaittimeout=30
numprocs=1
stdout_logfile=/tmp/supervisord_worker_submit.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"

[program:worker_monitor]
command=python3 -u worker_monitor.py
directory=/home/runner/workspace
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaittimeout=30
numprocs=1
stdout_logfile=/tmp/supervisord_worker_monitor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONIOENCODING="utf-8"
